# cron/com.phez.cron.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.phez.cron</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/caseymeehan/Documents/phez/cron/runner.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/caseymeehan/Documents/phez</string>

    <key>StartInterval</key>
    <integer>600</integer>

    <key>StandardOutPath</key>
    <string>/Users/caseymeehan/Documents/phez/cron/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/caseymeehan/Documents/phez/cron/logs/launchd-stderr.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
